# gmail-mcp-agent.service
[Unit]
Description=Gmail MCP Agent
After=network.target
Wants=network.target

[Service]
Type=simple
# Replace USER and the WorkingDirectory/paths below to match your deployment.
User=YOUR_USER
Group=YOUR_USER
WorkingDirectory=/opt/gmail-mcp-agent
ExecStart=/usr/bin/python3 /opt/gmail-mcp-agent/mcp_server.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PYTHONPATH=/opt/gmail-mcp-agent
Environment=PYTHONUNBUFFERED=1

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/gmail-mcp-agent

[Install]
WantedBy=multi-user.target
